# infra/systemd/fetchcat-n8n-backup.service
[Unit]
Description=Back up FetchCat n8n PostgreSQL and state
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
User=hanna
Group=hanna
ExecStart=/data/n8n/scripts/backup.sh
Nice=10

